# server/ingresso-uscita.service
# ======================================================================
# Servizio systemd per Server Ingresso/Uscita
# ======================================================================
#
# INSTALLAZIONE MANUALE:
# 1. Modifica questo file sostituendo:
#    - YOUR_USERNAME con il tuo username Linux
#    - /home/YOUR_USERNAME/ingresso_uscita_server con il percorso corretto
#
# 2. Copia il file in /etc/systemd/system/:
#    sudo cp ingresso-uscita.service /etc/systemd/system/
#
# 3. Ricarica systemd:
#    sudo systemctl daemon-reload
#
# 4. Abilita il servizio per l'avvio automatico:
#    sudo systemctl enable ingresso-uscita
#
# 5. Avvia il servizio:
#    sudo systemctl start ingresso-uscita
#
# 6. Verifica lo stato:
#    sudo systemctl status ingresso-uscita
#
# ======================================================================
# COMANDI UTILI:
# ======================================================================
#
# Avviare il servizio:
#   sudo systemctl start ingresso-uscita
#
# Fermare il servizio:
#   sudo systemctl stop ingresso-uscita
#
# Riavviare il servizio:
#   sudo systemctl restart ingresso-uscita
#
# Stato del servizio:
#   sudo systemctl status ingresso-uscita
#
# Abilitare avvio automatico:
#   sudo systemctl enable ingresso-uscita
#
# Disabilitare avvio automatico:
#   sudo systemctl disable ingresso-uscita
#
# Visualizzare i log:
#   sudo journalctl -u ingresso-uscita -f          # Log in tempo reale
#   sudo journalctl -t ingresso-uscita -f          # Log filtrati per tag
#   sudo journalctl -u ingresso-uscita -n 100      # Ultimi 100 log
#   sudo journalctl -u ingresso-uscita --since today  # Log di oggi
#
# ======================================================================

[Unit]
Description=Server Ingresso/Uscita - Sistema Gestione Presenze
Documentation=https://github.com/fragarray/ingresso_uscita
After=network.target
Wants=network-online.target

[Service]
Type=simple

# Sostituisci con il tuo username
User=YOUR_USERNAME

# Sostituisci con il percorso corretto del progetto
WorkingDirectory=/home/YOUR_USERNAME/ingresso_uscita_server

# Comando per avviare il server (rileva automaticamente node)
ExecStart=/usr/bin/env node /home/YOUR_USERNAME/ingresso_uscita_server/server.js

# Riavvio automatico (always = riavvia sempre, on-failure = solo se crasha)
Restart=always
RestartSec=10

# Numero massimo di riavvii in un periodo
StartLimitInterval=200
StartLimitBurst=5

# Log tramite syslog (più pulito e integrato con systemd)
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ingresso-uscita

# Variabili d'ambiente
Environment=NODE_ENV=production
Environment=PORT=3000

# Limiti di risorse (opzionali)
LimitNOFILE=4096

# Sicurezza (opzionali)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
